# files/flatpak-preinstall/flatpak-preinstall.service
[Unit]
Description=Pre-install initial Flatpak applications
Documentation=https://github.com/starlit-os/krytis/issues/66
# Retry until successful (marker absent means not yet done).
ConditionPathExists=!/var/lib/flatpak/.krytis-preinstall-done
Wants=network-online.target flatpak-system-helper.service
After=network-online.target flatpak-system-helper.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/krytis/flatpak-preinstall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
