# files/starlit-update/starlit-update-user.timer
[Unit]
Description=Periodic mise tool update
Documentation=https://github.com/starlit-os/krytis/issues/173

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
