# quadlet/buildbarn/bb-storage.container
[Unit]
Description=krytis Buildbarn storage (CAS + AC + FSAC)

[Container]
Image=ghcr.io/buildbarn/bb-storage@sha256:155fa806de0933e6222681353e2e9fe2f51bd2a78f88dd254ae0b5989fb8e585
ContainerName=bb-storage
Network=buildbarn.network
Exec=-config /config/storage.jsonnet

# gRPC (CAS/AC/FSAC) — mTLS terminated in-process, see config/storage.jsonnet.
PublishPort=7982:8981
# Prometheus/pprof diagnostics — bind to localhost only, not for CI/dev consumption.
PublishPort=127.0.0.1:9981:9981

Volume=%h/.local/share/buildbarn/config:/config:ro,Z
Volume=%h/.local/share/buildbarn/certs:/certs:ro,Z
Volume=buildbarn-storage-cas.volume:/data/storage-cas
Volume=buildbarn-storage-ac.volume:/data/storage-ac
Volume=buildbarn-storage-fsac.volume:/data/storage-fsac

[Service]
Restart=always

[Install]
WantedBy=default.target
